# genzquotes.service
[Unit]
Description=Gen Z Quotes Translator
After=network.target

[Service]
Type=simple
User=admin
Group=admin
WorkingDirectory=/var/www/quotes.iotacloud.tech
ExecStart=/usr/bin/node server/app.js
Restart=always
RestartSec=10
Environment=NODE_ENV=production
Environment=PORT=3000
Environment=OPENAI_API_KEY=your_openai_api_key_here

# Basic security settings
NoNewPrivileges=true
PrivateTmp=true

[Install]
WantedBy=multi-user.target 